# nginx-proxy-manager-mcp/rootfs/etc/services.d/npm-mcp/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "Nginx Proxy Manager MCP server stopped."
