# nginx-proxy-manager-mcp/rootfs/etc/services.d/npm-mcp/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Configure and start the Nginx Proxy Manager MCP server
# ==============================================================================

bashio::log.info "Configuring Nginx Proxy Manager MCP Add-on v$(bashio::addon.version)..."

export NPM_API_URL=$(bashio::config 'npm_api_url')
export NPM_IDENTITY=$(bashio::config 'npm_identity')
export NPM_SECRET=$(bashio::config 'npm_secret')
export NPM_MCP_HOST=$(bashio::config 'mcp_host')
export NPM_MCP_PORT=$(bashio::config 'mcp_port')
export NPM_MCP_TRANSPORT=$(bashio::config 'mcp_transport')

if bashio::config.has_value 'npm_proxy_defaults'; then
    export NPM_PROXY_DEFAULTS=$(bashio::config 'npm_proxy_defaults')
fi

bashio::log.info "NPM API: ${NPM_API_URL}"
bashio::log.info "MCP: ${NPM_MCP_TRANSPORT} on ${NPM_MCP_HOST}:${NPM_MCP_PORT}"

exec python3 -m npm_mcp.main
