# nginx-proxy-manager-mcp/rootfs/etc/services.d/npm-mcp/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Nginx Proxy Manager MCP server
# ==============================================================================

NPM_API_URL=$(bashio::config 'npm_api_url')
NPM_IDENTITY=$(bashio::config 'npm_identity')
NPM_SECRET=$(bashio::config 'npm_secret')
MCP_HOST=$(bashio::config 'mcp_host')
MCP_PORT=$(bashio::config 'mcp_port')

NPM_PROXY_DEFAULTS=""
if bashio::config.has_value 'npm_proxy_defaults'; then
    NPM_PROXY_DEFAULTS=$(bashio::config 'npm_proxy_defaults')
fi

bashio::log.info "Starting Nginx Proxy Manager MCP Server..."
bashio::log.info "NPM API URL: ${NPM_API_URL}"
bashio::log.info "MCP listening on ${MCP_HOST}:${MCP_PORT}"

export NPM_API_URL="${NPM_API_URL}"
export NPM_IDENTITY="${NPM_IDENTITY}"
export NPM_SECRET="${NPM_SECRET}"
export NPM_MCP_HOST="${MCP_HOST}"
export NPM_MCP_PORT="${MCP_PORT}"
export NPM_MCP_TRANSPORT="http"

if [ -n "${NPM_PROXY_DEFAULTS}" ]; then
    export NPM_PROXY_DEFAULTS="${NPM_PROXY_DEFAULTS}"
fi

exec python3 -m npm_mcp.main
