# nginx-proxy-manager-mcp/rootfs/etc/s6-overlay/s6-rc.d/npm-mcp/finish
#!/usr/bin/execlineb -S0
foreground { s6-echo "npm-mcp service stopped" }
if { s6-test ${1} -ne 0 }
foreground { s6-echo "npm-mcp exited with error ${1}" }
